# etc/profile.d/cc-hpc-shortcuts_1.0_all.sh
export SSM_SHORTCUT_PATH="$(cc-hpc-shortcuts_1.0_all-cfg)/local $SSM_SHORTCUT_PATH"
